# templates/vector.service.j2
# templates/vector.service.j2
[Unit]
Description=Vector Log Aggregator
Documentation=https://vector.dev
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User={{ vector_user }}
Group={{ vector_group }}
ExecStart={{ vector_path }}/vector --config {{ vector_config }}
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target